# provisioning/templates/etc/profile.d/hadoop.sh.j2
export HADOOP_PREFIX={{hadoop_install_dir}}