# provisioning/templates/etc/profile.d/hadoop.sh.j2
export HADOOP_PREFIX=/usr/local/lib/hadoop-2.7.1